# services/curb.service
[Unit]
Description=Service file for Curby Fan Controller

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/curby
